# ROOT/etc/profile.d/zzz-mqm.sh

export MQM=/opt/mqm
export PATH="$MQM/bin:$PATH"

