# ROOT/etc/profile.d/zzz-informix.sh
: choose the default engine for login sessions

. /etc/informix/profile.d/suzyq

