# ROOT/etc/informix/profile.d/suzyq
: setup shell variables to directly access suzyq engine

set -a

# remove old informix bin
[[ -n "$INFORMIXDIR" ]] && PATH="${PATH/"$INFORMIXDIR/bin:"}"
INFORMIXDIR=/opt/informix-11.70.FC4
INFORMIXSERVER=suzyq_tcp
ONCONFIG=onconfig.suzyq
INFORMIXSQLHOSTS=/etc/informix/config.d/sqlhosts
PATH=$INFORMIXDIR/bin:$PATH

DB_LOCALE=en_us.8859-1
CLIENT_LOCALE=en_us.8859-1
SERVER_LOCALE=en_us.8859-1

DBDATE=dmy4/
DBCENTURY=C

set +a

